# deployment/tongcove.service
[Unit]
Description=Tongcove Django Application - Learn Chinese Perth
Documentation=https://github.com/your-repo/tongcove
After=network.target
# Removed postgresql.service dependency since using SQLite3
# Wants=redis.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/tongcove
Environment=DJANGO_SETTINGS_MODULE=deployment.production_settings
Environment=PYTHONPATH=/var/www/tongcove
Environment=PYTHONUNBUFFERED=1

# 主进程命令 (Main process command)
ExecStart=/var/www/tongcove/.venv/bin/gunicorn \
    --name tongcove \
    --workers 3 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 30 \
    --keep-alive 2 \
    --bind 127.0.0.1:8000 \
    --pid /run/tongcove/tongcove.pid \
    --access-logfile /var/www/tongcove/logs/gunicorn_access.log \
    --error-logfile /var/www/tongcove/logs/gunicorn_error.log \
    --log-level info \
    tongcove.wsgi:application

# 进程管理 (Process management)
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=always
RestartSec=10

# 运行时目录 (Runtime directory)
RuntimeDirectory=tongcove
RuntimeDirectoryMode=0755

# 安全设置 (Security settings)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/tongcove
ReadWritePaths=/var/www/tongcove/logs
ReadWritePaths=/var/www/tongcove/media
ReadWritePaths=/run/tongcove

# 资源限制 (Resource limits)
LimitNOFILE=65536
LimitNPROC=4096

# 环境变量 (Environment variables)
EnvironmentFile=-/var/www/tongcove/.env.production

[Install]
WantedBy=multi-user.target